# templates/etc/systemd/system/shibboleth-idp-watcher.path
[Path]
PathModified=/etc/alternatives/jre/lib

[Install]
WantedBy=multi-user.target
