# templates/shibboleth-idp-watcher.path
[Path]
PathModified=/usr/lib/java

[Install]
WantedBy=multi-user.target
