# templates/shibboleth-idp-watcher.path
[Path]
PathModified=/etc/alternatives/java

[Install]
WantedBy=multi-user.target
